# deploy/systemd/numa-youtube-encoder.service
[Unit]
Description=Numa Radio — YouTube 24/7 broadcast encoder
After=network-online.target numa-liquidsoap.service
Wants=network-online.target
# Crash-loop protection — if the encoder fails 5 times in 5 minutes,
# stop trying. An operator should investigate before letting it spin
# (e.g. revoked stream key, network outage, broken /live page).
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=marku
EnvironmentFile=/etc/numa/env
# Belt-and-braces: kill any zombies from a previous run before starting.
# `[X]vfb` is a regex that matches "Xvfb" but the literal "[X]vfb" in
# the shell's own argv does not match the regex — without that trick
# pkill -f would SIGTERM its own /bin/sh and fail the whole start.
# We also kill orphan ffmpegs that match our x11grab :99 input — a
# leftover smoke-test or unclean restart can otherwise sit forever
# stealing Xvfb frames and CPU, which makes the real encoder choppy
# enough that YouTube buffers minutes ahead of live.
ExecStartPre=/bin/sh -c 'pkill -f "[X]vfb :99" || true; pkill -f "[g]oogle-chrome.*--display=:99" || true; pkill -f "[c]hromium.*--display=:99" || true; pkill -f "[f]fmpeg.*x11grab.*-i :99" || true; rm -f /tmp/numa-encoder-smoke.flv'
# WSL2 (WSLg) bind-mounts /tmp/.X11-unix as a read-only tmpfs, so Xvfb's
# unix socket creation fails. Replace it with a writable, sticky-bit
# directory on every start. `+` runs as root despite User=marku.
ExecStartPre=+/bin/sh -c 'umount /tmp/.X11-unix 2>/dev/null; mkdir -p /tmp/.X11-unix; chmod 1777 /tmp/.X11-unix'
ExecStart=/usr/local/bin/numa-youtube-encoder.sh
# Belt-and-braces cleanup. systemd's default KillMode=control-group
# already SIGTERMs every child in the cgroup, but the script does
# `exec ffmpeg` which replaces the bash trap. Run the same kill
# patterns as ExecStartPre so a stop leaves no stragglers.
ExecStopPost=/bin/sh -c 'pkill -f "[X]vfb :99" || true; pkill -f "[g]oogle-chrome.*--display=:99" || true; pkill -f "[c]hromium.*--display=:99" || true; pkill -f "[f]fmpeg.*x11grab.*-i :99" || true'
Restart=on-failure
RestartSec=15s
# Resource caps — Chromium can leak memory over weeks. Force a clean
# restart once a day so we always run from a fresh Chromium process.
RuntimeMaxSec=86400

[Install]
WantedBy=multi-user.target
