# deploy/systemd/numa-song-worker.service
[Unit]
Description=Numa Radio — song generation worker
After=network-online.target numa-queue-daemon.service
Wants=network-online.target
PartOf=multi-user.target

[Service]
Type=simple
User=marku
WorkingDirectory=/home/marku/saas/numaradio
# Production secrets live in /etc/numa/env (root-owned, mode 0600).
# .env.local is for local dev; production should fail loud if
# /etc/numa/env is missing rather than silently falling back to a
# committed dev-secret file. (No leading `-` = required.)
EnvironmentFile=/etc/numa/env
ExecStart=/usr/bin/npx tsx workers/song-worker/index.ts
Restart=on-failure
RestartSec=5
# Crash-loop protection — a recurring startup bug shouldn't fork an
# infinite restart loop into journald.
StartLimitBurst=3
StartLimitIntervalSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=numa-song-worker

KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
